# src/firmware-tools/unipifwupdate.service.in
[Unit]
Description=Unipi PLC Firmware Update Service
Before=unipitcp.service
ConditionVirtualization=No

[Service]
Type=oneshot
ExecStart=/bin/sh -c '@bindir@/fwspi --auto'
ExecStartPost=/bin/systemctl disable unipifwupdate

[Install]
WantedBy=multi-user.target
