# debian/systemd/neuron.target
[Unit]
Description=Neuron Target
After=basic.target
