# unipi-common/systemd/system/unipigprs.target

[Unit]
Description=Unipi GPRS Series Target
Before=unipispi.target
