# unipi-common/systemd/system/unipiS227.target

[Unit]
Description=Unipi S227 Series Target
Before=unipispi.target
