# unipi-common/systemd/system/unipiiqrf.target

[Unit]
Description=Unipi IQRF Series Target
#Conflicts=neuron.target
Before=unipispi.target
