# unipi-common/systemd/system/unipi-ttymxc2-rs485.target
[Unit]
Description=Unipi RS485 on ttymxc2 Target
#Conflicts=neuron.target
Before=unipispi.target
