# unipi-common/systemd/system/unipigate.target
[Unit]
Description=Gate Target
After=basic.target
