# unipi-common/systemd/system/unipinb.target

[Unit]
Description=Unipi NB Series Target
Before=unipispi.target
