# unipi-common/systemd/system/unipi-ttys0-rs485.target
[Unit]
Description=Unipi RS485 on ttyS0 Target
#Conflicts=neuron.target
Before=unipispi.target
